# host/master/scripts/systemd/agv-planner.service
[Unit]
Description=AGV Task Planner — reads SHM, plans AGV routes, dispatches via MQ
Documentation=https://github.com/openclaw/ENavigate
After=agv-model-manager.service agv-router.service agv-log-daemon.service
Requires=agv-model-manager.service agv-router.service agv-log-daemon.service

[Service]
Type=simple
User=root
Group=root
ExecStart=/opt/enavigate-master/build/src/planner/planner
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agv-planner
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
